# chat-reminder.service
[Unit]
Description=Chat Screenshot Reminder Web Service
After=network.target

[Service]
User=bytedance
WorkingDirectory=/Users/bytedance/Documents/trae_projects/client_chat_reminder
Environment="PATH=/Users/bytedance/Documents/trae_projects/client_chat_reminder/venv/bin"
EnvironmentFile=/Users/bytedance/Documents/trae_projects/client_chat_reminder/.env
ExecStart=/Users/bytedance/Documents/trae_projects/client_chat_reminder/venv/bin/uvicorn src.main:app --host 0.0.0.0 --port ${PORT:-80}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
